# assets/slapd.conf
[program:slapd]
autostart=true
user=openldap
group=openldap
command=/usr/sbin/slapd -h ldap:/// -h ldapi:/// -F /etc/ldap/slapd.d -d 1
